# deploy/systemd/pi-ups-dashboard-update.service
[Unit]
Description=Pi UPS Dashboard release updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/repos/pi-ups-dashboard
EnvironmentFile=-/etc/default/pi-ups-dashboard-updater
ExecStart=/opt/repos/pi-ups-dashboard/deploy/systemd/update-from-release.sh

[Install]
WantedBy=multi-user.target
